# src/cordova/build/modules/build/files/etc/supervisor/conf.d/consul.conf
[program:consul]
command = /usr/local/bin/consul agent -client 0.0.0.0 -config-dir /etc/consul/consul.d -data-dir /consul/data
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
